# radiod.service
[Unit]
Description=radiod - Internet Radio Daemon
After=network-online.target sound.target
Wants=network-online.target sound.target

[Service]
Type=simple
ExecStart=radiod
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/radiod
ReadWritePaths=%h/.cache
ReadOnlyPaths=/usr/lib
ReadOnlyPaths=/usr/share

[Install]
WantedBy=default.target
